# gcs/newfiles_skel/lib/udev/rules.d/40-gecos_usb_modeswitch.rules
# Guadalinex GECOS extra usb_modeswitch rules

ATTRS{idVendor}=="12d1", ATTRS{idProduct}=="1f1c", RUN+="usb_modeswitch '%b/%k'"
